# binance-service.service
[Unit]
Description=Binance Trade Aggregation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/binance_service /etc/binance-service/config.json
WorkingDirectory=/etc/binance-service
Restart=always
RestartSec=5
User=binance-service
Group=binance-service

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/binance-service

[Install]
WantedBy=multi-user.target